# templates/usr/lib/systemd/system/carbon-c-relay.service.erb
# This file is managed by Puppet. DO NOT EDIT.
[Unit]
Description=Enhanced C implementation of Carbon relay, aggregator and rewriter
After=network.target

[Service]
User=<%= @user %>
Group=<%= @group %>
PIDFile=<%= @pid_dir %>/<%= @pid_file %>
EnvironmentFile=/etc/sysconfig/<%= @service_name %>
ExecStart=/usr/bin/carbon-c-relay -P <%= @pid_dir %>/<%= @pid_file %> -f <%= @config_file %> ${ARGS}
ExecReload=/bin/kill -HUP ${MAINPID}
Type=simple
Restart=on-failure
<% if @limit_fsize -%>
# (file size)
LimitFSIZE=<%= @limit_fsize %>
<% end -%>
<% if @limit_cpu -%>
# (cpu time)
LimitCPU=<%= @limit_cpu %>
<% end -%>
<% if @limit_as -%>
# (virtual memory size)
LimitAS=<%= @limit_as %>
<% end -%>
# (open files)
LimitNOFILE=<%= @limit_nofile %>
# (processes/threads)
LimitNPROC=<%= @limit_nproc %>

[Install]
WantedBy=multi-user.target
